# luci-app-wolplus/root/etc/uci-defaults/luci-app-wolplus
#!/bin/sh

rm -f /tmp/luci-indexcache /tmp/luci-modulecache/*
exit 0
